# src/variants/SimplyPi/filesystem/etc/systemd/system/simplypi.service
[Unit]
Description=SimplyPrint startup service

[Service]
Type=oneshot
ExecStart=bash /home/root/bin/simplyprint_start.sh

[Install]
WantedBy=multi-user.target
