# src/variants/SimplyPi/filesystem/root/etc/systemd/system/simplypi.service
[Unit]
Description=SimplyPrint startup service
[Service]
ExecStart=/usr/bin/simplyprint_start
[Install]
WantedBy=multi-user.target
